# modules/nodepool/files/nodedrain.service
[Unit]
Description=node drain
After=k3s.service systemd-resolved.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/echo "Node drain reporting for duty"
TimeoutStopSec=120s
ExecStop=/bin/bash /usr/local/bin/nodedrain.sh

[Install]
WantedBy=multi-user.target
